# amdgpu-fancontrol.service
[Service]
Type=simple
ExecStart=/usr/bin/amdgpu-fancontrol
